# manage_modem.service
[Unit]
Description=Modem Manager - Centralized modem control service
After=network.target
Wants=network.target
After=network.target asterisk.service
Requires=asterisk.service
# Wait for device to be available
After=dev-ttyUSB2.device
Requires=dev-ttyUSB2.device

[Service]
Type=simple
ExecStart=/usr/bin/python3 /mnt/data/manage_modem.py \
    --port 5555 \
    --host 0.0.0.0 \
    --modem /dev/ttyUSB2 \
    --baud 115200 \
    --log-level INFO \
    --log-file /mnt/data/calls.log
WorkingDirectory=/mnt/data
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
User=root

# Resource limits
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
