# manage_modem.service
[Unit]
Description=Modem Manager - Centralized modem control service
After=network.target
Wants=network.target
After=network.target
# Wait for device to be available
After=dev-ttyUSB2.device
Requires=dev-ttyUSB2.device

[Service]
Type=simple
User=kuser
ExecStartPre=/bin/sleep 3
# Use explicit UID for kuser
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="PULSE_RUNTIME_PATH=/run/user/1000/pulse"
Environment="PULSE_SERVER=unix:/run/user/1000/pulse/native"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus"
# Debug: Print environment variables
ExecStartPre=/bin/sh -c 'echo "[MANAGE_MODEM_SERVICE] XDG_RUNTIME_DIR=$XDG_RUNTIME_DIR"'
ExecStartPre=/bin/sh -c 'echo "[MANAGE_MODEM_SERVICE] PULSE_RUNTIME_PATH=$PULSE_RUNTIME_PATH"'
ExecStartPre=/bin/sh -c 'echo "[MANAGE_MODEM_SERVICE] PULSE_SERVER=$PULSE_SERVER"'
ExecStartPre=/bin/sh -c 'echo "[MANAGE_MODEM_SERVICE] DBUS_SESSION_BUS_ADDRESS=$DBUS_SESSION_BUS_ADDRESS"'
ExecStart=/usr/bin/python3 /mnt/data/manage_modem.py \
    --port 5555 \
    --host 0.0.0.0 \
    --modem /dev/ttyUSB2 \
    --baud 115200 \
    --log-level INFO \
    --log-file /mnt/data/calls.log
WorkingDirectory=/mnt/data
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal


# Resource limits
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
