# VOIP/voip_call_monitor.service
[Unit]
Description=VoIP Call Monitor TCP Service
After=network.target asterisk.service manage_modem.service
Requires=asterisk.service
# Wait for device to be available
After=dev-ttyUSB2.device
Requires=dev-ttyUSB2.device

[Service]
Type=simple
User=kuser
Group=kuser
WorkingDirectory=/mnt/data
ExecStart=/usr/bin/python3 /mnt/data/voip_call_monitor_tcp.py -r
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
# Kill only the main process on stop
KillMode=mixed
# Give the process 10 seconds to stop gracefully before killing
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
